# com.jobbot.slackbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobbot.slackbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/nicoleeid/jobbot_clean/start_working_bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nicoleeid/jobbot_clean</string>

    <key>StandardOutPath</key>
    <string>/Users/nicoleeid/jobbot_clean/slack_bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nicoleeid/jobbot_clean/slack_bot_error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/nicoleeid</string>
    </dict>
</dict>
</plist>
